# ext/debian/mcollective.postrm
#!/bin/bash -e

#####
# init script deregistration
update-rc.d -f mcollective remove &> /dev/null || true

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
